# scripts/systemd/mirrorcrm-cleanup.timer
[Unit]
Description=Daily MirrorCRM files cleanup timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
